# agent-hub/com.agent-hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-hub</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/boty/nanoclaw/agent-hub/venv/bin/python3.13</string>
        <string>orchestrator.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/boty/nanoclaw/agent-hub</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/boty/.local/bin</string>
        <key>HOME</key>
        <string>/Users/boty</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/boty/nanoclaw/agent-hub/logs/agent-hub.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/boty/nanoclaw/agent-hub/logs/agent-hub.error.log</string>
</dict>
</plist>
